# docker/s6-rc.d/mcp-bridge-http/finish
#!/bin/sh
echo "MCP Bridge HTTP stopped with exit code $1"
